# pi_motor_bridge/motor-bridge.service
[Unit]
Description=Sirena Nina motor bridge (RPi side)
After=multi-user.target pigpiod.service
Wants=pigpiod.service

[Service]
Type=simple
# Run as root because pigpio needs GPIO access. If you've set up
# pigpio for an unprivileged user, change User= to that user and
# adjust the working dir / install path accordingly.
User=root
ExecStart=/usr/bin/python3 /opt/sirena/pi_motor_bridge/motor_bridge.py
Restart=on-failure
RestartSec=2
# Park the motors on stop / restart.
KillSignal=SIGTERM
TimeoutStopSec=5

# Optional environment overrides (uncomment + edit):
# Environment=NINA_BRIDGE_PORT=/dev/serial0
# Environment=NINA_BRIDGE_BAUD=115200
#
# Watchdog: stops wheels if no command for N seconds while moving.
# Default 1.5s. Set to 0 to DISABLE the watchdog entirely - useful
# when debugging erratic stops / direction glitches and you want the
# wheels to keep going until the Jetson explicitly says STOP. Re-enable
# by setting back to 1.5 (or commenting this line out).
# Environment=NINA_BRIDGE_WATCHDOG_SEC=0

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
